# panasonic-brightness.service
[Unit]
Description=Panasonic Brightness
[Service]
ExecStart=/usr/local/bin/brightness
[Install]
WantedBy=basic.target